# crontab
# Scraper-Pro Cron Jobs
# Validation: process pending contacts every hour
0 * * * * . /etc/environment; cd /app && python -m scraper.jobs.process_contacts >> /app/logs/process.log 2>&1

# MailWizz sync: inject validated contacts every hour (offset 30 min)
30 * * * * . /etc/environment; cd /app && python -m scraper.jobs.sync_to_mailwizz >> /app/logs/mailwizz_sync.log 2>&1

